# mycustom.slice
[Unit]
Description=My Custom Slice

[Slice]
CPUQuota=1%
MemoryMax=1M
